# host/systemd/vps-proxy.service
[Unit]
Description=VPS TLS SNI Passthrough Proxy (port 443)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/vps-mcp/host.env
ExecStart=/usr/local/sbin/vps-proxy -domain ${DOMAIN}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
