# host/systemd/vps-proxy80.socket
[Unit]
Description=VPS HTTP TCP Passthrough Proxy socket (port 80)
Documentation=https://github.com/paijp/vps-subdomain-mcp

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
