# host/systemd/vps-proxy443.socket
[Unit]
Description=VPS TLS SNI Passthrough Proxy socket (port 443)
Documentation=https://github.com/paijp/vps-subdomain-mcp

[Socket]
ListenStream=443
Accept=no

[Install]
WantedBy=sockets.target
