# host/systemd/vps-proxy80.socket
[Unit]
Description=VPS HTTP TCP Passthrough Proxy socket (port 80)

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
