# host/systemd/vps-proxy-http.service
[Unit]
Description=VPS HTTP Reverse Proxy (port 80)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/vps-mcp/host.env
ExecStart=/usr/local/sbin/vps-proxy-http -domain ${DOMAIN}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
